# usr/lib/systemd/system/ffmpeg-camera.service
[Unit]
Description=ffmpeg-camera

[Service]
ExecStart=/bin/sh /root/website/run2.sh

[Install]
WantedBy=multi-user.target
