# usr/lib/systemd/system/ipstatic.service
[Unit]
Description=ipstatic

[Service]
ExecStart=/bin/bash /root/beagle/startup_server.sh

[Install]
WantedBy=multi-user.target
