# usr/lib/systemd/system/server.service
[Unit]
Description=server

[Service]
ExecStart= /usr/bin/node /root/beagle/server.js

[Install]
WantedBy=multi-user.target
