# usr/lib/systemd/system/server-stream.service
[Unit]
Description=server-stream

[Service]
ExecStart=/usr/bin/node /root/website/stream-server.js
[Install]
WantedBy=multi-user.target
