# usr/lib/systemd/system/stream.service
[Unit]
Description=stream

[Service]
ExecStart=/bin/bash /root/beagle/boneCV/streamVideoUDP

[Install]
WantedBy=multi-user.target
